# example/crontab
0 8 * * Mon php /var/www/vendor/alsvanzelf/debby/notify.php example/project personal-access-token
0 8 * * Mon php /var/www/vendor/alsvanzelf/debby/notify.php /var/www/debby-options.json
