# example/crontab
0 8 * * Mon export php /var/www/vendor/alsvanzelf/debby/notify.php devops@example.com
0 8 * * Mon export php /var/www/vendor/alsvanzelf/debby/notify.php /var/www/debby-options.json
